# supervisor/conf.d/confd.conf
[program:confd]
command=/usr/local/bin/confd -node=%(ENV_ETCD_ENDPOINT)s -interval=10
autostart=true
autorestart=true
user=root
startsecs=6
priority=999