# supervisor/conf.d/confd.conf
[program:confd]
command=/usr/local/bin/confd -node $ETCD_ENDPOINT -interval 10 -config-file /etc/confd/conf.d/haproxy.toml
autostart=true
autorestart=true
user=root
startsecs=6
priority=999