# supervisor/conf.d/haproxy.conf
[program:haproxy]
command=/usr/local/bin/run-haproxy.sh
autostart=true
autorestart=true
user=root
startsecs=6
priority=999